# var/spool/cron/crontabs/edge
* * * * * php /var/www/bin/magento cron:run

