# var/spool/cron/crontabs/edge
MAILTO=""
* * * * * php /var/www/bin/magento cron:run 2>&1

