# crontab
*/5 * * * * www-data /bin/sh /var/www/cron.sh