# etc/crontab
PATH=/usr/local/bin:/bin:/usr/bin

*/5 * * * * www-data /var/www/cron.sh 2>&1
0 4 * * * www-data php -f /var/www/shell/indexer.php reindexall 2>&1
