# etc/crontab
PATH=/usr/local/bin:/bin:/usr/bin

*/5 * * * * www-data /var/www/cron.sh 2>&1