# etc/crontab
PATH=/usr/local/bin:/bin:/usr/bin
MAILTO=""

*/5 * * * * www-data /bin/sh /var/www/cron.sh
0 4 * * * www-data php -f /var/www/shell/indexer.php reindexall

